# systemd./user/mopidy.service
[Unit]
Description=mopidy
After=network.target

[Service]
ExecStart=/usr/bin/mopidy

[Install]
WantedBy=multi-user.target
