# systemd./user/udiskie.service
[Unit]
Description=mopidy
After=network.target

[Service]
ExecStart=/usr/bin/udiskie

[Install]
WantedBy=multi-user.target
